# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1280x900x24 -ac
autostart=true
autorestart=false
startsecs=0
priority=1

[program:fluxbox]
command=/usr/bin/fluxbox
autostart=false
autorestart=false
environment=DISPLAY=:0
priority=2

[program:x11vnc]
command=/usr/bin/x11vnc -display :0 -forever -nopw
autostart=false
autorestart=false
environment=DISPLAY=:0
priority=3

[program:novnc]
command=bash /root/novnc/utils/launch.sh --vnc localhost:5900 --listen 8080
autostart=false
autorestart=false
priority=4

[program:fastapi]
command=python -m uvicorn server:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
environment=DISPLAY=:0
priority=5
stdout_logfile=/var/log/supervisor/fastapi-stdout.log
stderr_logfile=/var/log/supervisor/fastapi-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=2MB
stdout_logfile_backups=10
stderr_logfile_backups=5
